# puppet/modules/librenmssimplebackup/files/librenms-simple-backup.cron
55 3 * * * root /usr/local/bin/librenms-simple-backup.sh >/dev/null 2>&1